# systemd/tiertune-aws.service
[Unit]
Description=Configure AWS instance
After=local-fs.target
ConditionPathExists=/etc/tiertune-aws.yml

[Service]
Type=oneshot
ExecStart=/usr/bin/tiertune-aws
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
